# supervisord/supervisord.conf
; Sample supervisor config file.

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:13148       ; (ip_address specifier, *:port for all iface)
username=super              ; (default is no username (open server))
password=abudabu1               ; (default is no password (open server))

[supervisord]
logfile=/home/entone/webapps/gericare/GeriCare/supervisord/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/home/entone/webapps/gericare/GeriCare/supervisord/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:13148 ; use an http:// url to specify an inet socket
username=super              ; should be same as http_username if set
password=abudabu1                ; should be same as http_password if set

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

[program:GeriCare]
command=python start.py              ; the program (relative uses PATH, can take args)
directory=/home/entone/webapps/gericare/GeriCare/web/                ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
exitcodes=0,2                 ; 'expected' exit codes for process (default 0,2)

[program:GeriCareAdmin]
command=python start.py              ; the program (relative uses PATH, can take args)
directory=/home/entone/webapps/gericare/GeriCare/admin/                ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
exitcodes=0,2                 ; 'expected' exit codes for process (default 0,2)

; The below sample group section shows all possible group values,
; create one or more 'real' group: sections to create "heterogeneous"
; process groups.

[group:GeriCare]
programs=GeriCare,GeriCareAdmin ; each refers to 'x' in [program:x] definitions
